# deploy/gce/aster-daily-restart.service
[Unit]
Description=Update + Restart Aster Strategy After Daily Maintenance Window

[Service]
Type=oneshot
ExecStart=/opt/aster/deploy/gce/daily_update_and_restart.sh
TimeoutStartSec=1800
