# deploy/gce/aster-daily-restart.service
[Unit]
Description=Restart Aster Strategy After Daily Maintenance Window

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart aster.service
