# deploy/gce/aster.service
[Unit]
Description=Aster Strategy Runtime
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=aster
Group=aster
WorkingDirectory=/opt/aster/core
EnvironmentFile=/opt/aster/deploy/gce/aster.env
ExecStart=/opt/aster/deploy/gce/run_strategy.sh
# Restart only on crashes/non-zero exits. If core/main.py exits cleanly after
# ASTER_POLL_TIME seconds, systemd will not restart it.
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGINT
NoNewPrivileges=true
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
